# ops/hermes.service.example
[Unit]
Description=Hermes email coding agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hermes
Group=hermes
WorkingDirectory=/srv/hermes
EnvironmentFile=/srv/hermes/.env
ExecStart=/srv/hermes/.venv/bin/python -m hermes
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/hermes

[Install]
WantedBy=multi-user.target
